# packages/telegram-remote/examples/systemd/skc-telegram-remote-rpc.service
# Linux systemd user unit for the Telegram Remote RPC backend.
# Copy to: ~/.config/systemd/user/skc-telegram-remote-rpc.service
# The bot dials an existing socket; it never owns or kills the RPC session.
[Unit]
Description=SKC Telegram Remote gateway (RPC backend)
Documentation=https://github.com/jaybeyond/Sayknow_CLI/tree/dev/packages/telegram-remote
Wants=network-online.target skc-rpc-session.service
After=network-online.target skc-rpc-session.service

[Service]
Type=simple
WorkingDirectory=%h/src/sayknow-cli/packages/telegram-remote
Environment=PATH=%h/.bun/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/.config/skc/telegram-remote.env
ExecStartPre=/bin/sh -lc 'socket="$SKC_TELEGRAM_REMOTE_RPC_SOCKET"; i=0; while [ "$i" -lt 30 ]; do [ -S "$socket" ] && exit 0; i=$((i+1)); sleep 1; done; echo "RPC socket not ready after 30s" >&2; exit 1'
ExecStart=/usr/bin/env bun run src/cli.ts
Restart=on-failure
RestartSec=10s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
