# packages/telegram-remote/examples/launchd/com.example.sayknow-cli.skc-rpc-session.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.example.sayknow-cli.skc-rpc-session</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/YOU/.local/bin/skc</string>
		<string>launch</string>
		<string>--output</string>
		<string>rpc</string>
		<string>--listen</string>
		<string>/Users/YOU/Library/Caches/skc/telegram-remote.sock</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/YOU/.bun/bin:/Users/YOU/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/YOU/src/sayknow-cli</string>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/YOU/Library/Logs/skc-rpc-session.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOU/Library/Logs/skc-rpc-session.err.log</string>
</dict>
</plist>
